# setup/packages/safecor/safecor-core/domu_apkovl.tmpl/etc/init.d/psec-core-init
#!/sbin/openrc-run

name="safecor-core-init"

depend() {
    after syslog
}

start() {
    ebegin "Finalizing Safecor initialization"

    apk update
    apk upgrade
    apk add

    eend $?
}

stop() {
    ebegin ""
    eend $?
}
