# setup/packages/psec/psec-diag/psec-diag.post-install
#!/bin/sh

rc-service udev start
rc-service udev-trigger start
rc-update add udev
rc-update add udev-trigger

rc-update add diag
rc-service diag start

exit 0