# setup/packages/psec/psec-diag/psec-diag.post-install
#!/bin/sh

rc-update add diag
rc-service diag start

exit 0