# setup/packages/psec/psec-diag-gui/psec-diag-gui.post-install
#!/bin/sh

rc-update add psec-diag
rc-service psec-diag start

exit 0