# setup/packages/psec/psec-diag-gui/start-gui.initd
#!/sbin/openrc-run

name="start-gui"
command_background="yes"

depend() {    
    after syslog
}

start() {
    ebegin "Starting GUI"

    chmod +x /usr/lib/psec/diag/src/diag/start.sh
    mkdir -p /var/log/psec
    /usr/lib/psec/diag/src/diag/start.sh &
}

stop() {
    ebegin "Nothing"
    
    eend $?
}