# setup/packages/safecor/safecor-core/init/debug.initd
#!/sbin/openrc-run

name="debug"

depend() {
    after mosquitto
}

start() {
    if grep -q "DEBUG=ON" "/proc/cmdline"; then            
        ebegin "Create TTY tunnel"
        start-stop-daemon --start --make-pidfile --exec /usr/lib/safecor/bin/create-tty-tunnel.sh --pidfile /var/run/$name-tty-tunnel.pid --background --user "svc-tty-tunnels"
        eend $?

        ebegin "Create TTY terminal"
        start-stop-daemon --start --make-pidfile --exec /usr/lib/safecor/bin/create-tty-terminal.sh --pidfile /var/run/$name-tty-terminal.pid --background --user "svc-tty-terminal"
        eend $?

        ebegin "Log MQTT messages"
        start-stop-daemon --start --make-pidfile --exec /usr/lib/safecor/bin/start-mqtt-logger.sh --pidfile /var/run/$name-mqtt.pid --background --user "svc-mqtt-logger"
        eend $?

        # Store the information into Xenstore
        xenstore-write /local/domain/system/debug_on 1
    else 
        echo "Debugging disabled"
        eend 1
    fi
}

stop() {
    ebegin "Destroy TTY tunnel"
    kill -9 $(cat /var/run/$name-tty-tunnel.pid)
    eend $?

    ebegin "Destroy TTY terminal"
    kill -9 $(cat /var/run/$name-tty-terminal.pid)
    eend $?

    ebegin "Stop MQTT logging"
    kill -9 $(cat /var/run/$name-mqtt.pid)
    eend $?
}