# setup/packages/safecor/safecor-tests-gui/psec-tests-gui.post-install
#!/bin/sh

mkdir -p /var/log/safecor

chmod +x /usr/lib/safecor/tests/src/tests/start.sh

rc-update add tests
rc-service tests start

# Remove unneeded big files
rm -rf /usr/lib/libQt6WebEngineCore*

exit 0