# setup/packages/psec/psec-sys-usb/99-usbdisks.rules
ENV{ID_FS_USAGE}=="filesystem|other|crypto", ENV{ID_FS_LABEL_ENC}=="?*", RUN+="/usr/local/psec/bin/mdev-usb-storage"