# setup/packages/psec/psec-sys-usb/99-usbdisks.rules
SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", ENV{ID_FS_USAGE}=="filesystem", RUN+="/usr/lib/psec/bin/mdev-usb-storage.sh"
