# setup/packages/psec/psec-sys-usb/psec-sys-usb-controller.initd
#!/sbin/openrc-run

name="psec-sys-usb-controller"
user="svc-sys-usb-controller"

depend() {
    before udev
}

start() {
    ebegin "Starting PSEC sys-usb controller"

    start-stop-daemon --start --make-pidfile --exec /usr/bin/python3 --pidfile /var/run/$name.pid --background --user "$user" -- /usr/lib/psec/bin/psec-sys-usb-controller.py

    eend $?
}

stop() {
    ebegin "Stopping PSEC sys-usb controller"

    kill -9 $(cat /var/run/$name.pid)

    eend $?
}