# setup/packages/psec/psec-sys-usb/psec-sys-usb.post-install
#!/bin/sh

mount -o remount,size=2M /dev/shm

mkdir -p /var/log/psec

modprobe exfat
modprobe vfat
modprobe ntfs3
modprobe ext4

rc-update add psec-sys-usb-controller
rc-service psec-sys-usb-controller start

echo Start udev service
#rc-service mdev stop
#rc-update del mdev sysinit
rc-update add udev boot
rc-service udev start
udevadm control --reload-rules

rc-update add tty-tunnel
rc-service tty-tunnel start

/usr/lib/psec/bin/hardening.sh

exit 0