# setup/packages/psec/psec-sys-usb/psec-sys-usb.post-install
#!/bin/sh

mkdir -p /var/log/psec

modprobe exfat
modprobe vfat

rc-service psec-sys-usb-controller start

echo Start/stop mdev/udev services
rc-service mdev stop
rc-update del mdev sysinit
rc-update add udev boot
rc-service udev start

exit 0