# setup/packages/psec/psec-sys-usb/psec-sys-usb.post-install
#!/bin/sh

mkdir -p /var/log/psec

modprobe exfat
modprobe vfat
modprobe ntfs

rc-update add psec-sys-usb-controller
rc-service psec-sys-usb-controller start

echo Start udev service
#rc-service mdev stop
#rc-update del mdev sysinit
rc-update add udev boot
rc-service udev start
udevadm control --reload-rules

exit 0