# setup/packages/psec/psec-sys-usb/psec-sys-usb.post-install
#!/bin/sh

rc-update add mdev sysinit
#rc-update add udev sysinit
rc-update add psec-sys-usb-controller