# setup/packages/psec/psec-sys-usb/tty-tunnel.initd
#!/sbin/openrc-run

name="tty-tunnel"
user="svc-tty-tunnel"

depend() {
    after udev
}

start() {
    ebegin "Create TTY tunnel(s)"

    start-stop-daemon --start --make-pidfile --exec /usr/lib/psec/bin/create-tty-tunnel.sh --pidfile /var/run/$name.pid --background --user "$user"

    eend $?
}

stop() {
    ebegin "Destroy TTY tunnel(s)"

    kill -9 $(cat /var/run/$name.pid)

    eend $?
}