# setup/packages/safecor/safecor-demo-gui/psec-demo-gui.post-install
#!/bin/sh

mkdir -p /var/log/safecor

chmod +x /usr/lib/safecor/diag/src/diag/start.sh

rc-update add start-gui
rc-service start-gui start

# Remove unneeded big files
rm -rf /usr/lib/libQt6WebEngineCore*

exit 0