# setup/packages/safecor/safecor-sys-usb/99-usbdisks.rules
SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", ENV{ID_FS_USAGE}=="filesystem", RUN+="/usr/lib/safecor/bin/mdev-usb-storage.sh"
