# setup/packages/safecor/safecor-core/psec-core.post-install
#!/bin/sh

#if [ "$APK_PHASE" != "post-install" ]; then
#    exit 0
#fi

echo
echo --------------------------------------------
echo -------      Safecor Installation       -------
echo ---        script post-install.sh        ---
echo --------------------------------------------
echo

hostname safecor.local
echo "" >> /etc/profile
echo "export DISPLAY=:0" >> /etc/profile

echo "**** Resizing shared memory"
# Note : non-persistent setting
mount -o remount,size=16M /dev/shm

echo "**** Creating service accounts"
addgroup safecor
adduser -H -D -s /sbin/nologin -G safecor svc-mqtt-tunnels
adduser -H -D -s /sbin/nologin -G safecor svc-tty-tunnels
adduser -H -D -s /sbin/nologin -G safecor svc-tty-terminal
adduser -H -D -s /sbin/nologin -G safecor svc-mqtt-logger
adduser -H -D -s /sbin/nologin -G safecor svc-orchestrator
adduser -H -D -s /sbin/nologin -G safecor svc-safecor-controller

echo "**** Settings permissions on files"
mkdir -p /var/log/safecor
chown -R nobody:safecor /var/log/safecor
chown -R nobody:safecor /usr/lib/safecor

rc-update add splash boot
rc-service splash start

rc-service modules restart

rc-service xenconsoled start
rc-update add xenconsoled
rc-service xenstored start
rc-update add xenstored
rc-update add mosquitto
rc-update add create-mqtt-tunnels
rc-update add safecor-core-controller
rc-update add start-x
rc-update add orchestrator
rc-update add debug
rc-update add hardening
rc-update add apparmor boot

# Install new Grub file
mv /etc/default/grub /etc/default/grub.orig
if [ -f /etc/default/grub.apk-new ]; then
    mv /etc/default/grub.apk-new /etc/default/grub
fi
# update-grub is automatically called at the end

# Re-enable colored prompt   
#mv "/etc/profile.d/color_prompt.sh.disabled" "/etc/profile.d/color_prompt.sh"

# Disable AppArmor unused profiles
mv /etc/apparmor.d/bin.ping /etc/apparmor.d/disable
mv /etc/apparmor.d/lsb_release /etc/apparmor.d/disable
mv /etc/apparmor.d/php-fpm /etc/apparmor.d/disable
mv /etc/apparmor.d/samba-* /etc/apparmor.d/disable
mv /etc/apparmor.d/sbin.klogd /etc/apparmor.d/disable
mv /etc/apparmor.d/sbin.syslog* /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.lib.apache2.mpm-prefork.apache2 /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.lib.dovecot.* /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.apache2 /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.avahi-daemon /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.dnsmasq /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.dovecot /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.identd /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.mdnsd /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.nmbd /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.nscd /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.ntpd /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.smbd /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.smbldap-useradd /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.traceroute /etc/apparmor.d/disable
mv /etc/apparmor.d/usr.sbin.winbindd /etc/apparmor.d/disable

echo Setup system integrity
/usr/lib/safecor/bin/integrity-setup.sh

exit 0
