# setup/packages/safecor/safecor-core/safecor-core.pre-install
#!/bin/sh

rc-service xenstored start
rc-update add xenstored