# setup/packages/psec/psec-demo-gui/psec-demo-gui.post-install
#!/bin/sh

mkdir -p /var/log/psec

chmod +x /usr/lib/psec/diag/src/diag/start.sh

rc-update add start-gui
rc-service start-gui start

# Remove unneeded big files
rm -rf /usr/lib/libQt6WebEngineCore*

exit 0