# setup/packages/psec/psec-core/init/splash.initd
#!/sbin/openrc-run

name="splash"

depend() {
    before modules
}

start() {
    width=`cat /sys/class/graphics/fb0/virtual_size | cut -d, -f1`
    height=`cat /sys/class/graphics/fb0/virtual_size | cut -d, -f2`
    splash_file="/boot/Splash_$width_$height.ppm"

    if [ ! -f $splash_file ]
    then 
        splash_file="/boot/Splash_1024_768.ppm"
    fi 

    ebegin "Loading splash file ($splash_file)"
    fbsplash -s $splash_file -T 2
    eend $?
}

stop() {
    ebegin ""
    eend $?
}
