# setup/packages/psec/psec-core/init/tty-terminal.initd
#!/sbin/openrc-run

name="tty-terminal"

depend() {
    after tty-tunnel
}

start() {
    ebegin "Create TTY terminal"

    start-stop-daemon --start --make-pidfile --exec /usr/lib/psec/bin/create-tty-terminal.sh --pidfile /var/run/$name.pid --background

    eend $?
}

stop() {
    ebegin "Destroy TTY terminal"

    kill -9 $(cat /var/run/$name.pid)

    eend $?
}