# setup/packages/psec/psec-core/init/tty-tunnel.initd
#!/sbin/openrc-run

name="tty-tunnel"

depend() {
    after create-mqtt-tunnels
}

start() {
    ebegin "Create TTY tunnel"

    if grep -q "DEBUG=ON" "/proc/cmdline"; then        
        start-stop-daemon --start --make-pidfile --exec /usr/lib/psec/bin/create-tty-tunnel.sh --pidfile /var/run/$name.pid --background
        eend $?
    else
        echo "Disabled"
        eend 1
    fi    
}

stop() {
    ebegin "Destroy TTY tunnel"

    kill -9 $(cat /var/run/$name.pid)

    eend $?
}