# setup/packages/psec/psec-syslog/psec-logger-daemon.initd
#!/sbin/openrc-run

name="psec-logger-daemon"
user="psec-logger-daemon"

depend() {
    before rsyslog
}

start() {         
    ebegin "Create logging daemon for rsyslog"
    start-stop-daemon --start --make-pidfile --exec /usr/bin/python3 --pidfile /var/run/$name.pid --background --user "$user" -- /usr/lib/psec/bin/psec-logger-daemon.py
    eend $?
}

stop() {
    ebegin "Destroy logging daemon for rsyslog"
    kill -9 $(cat /var/run/$name.pid)
    eend $?
}