# setup/packages/psec/psec-core/domu_apkovl.tmpl/etc/init.d/psec-core-init
#!/sbin/openrc-run

name="psec-core-init"

depend() {
    after syslog
}

start() {
    ebegin "Finalizing PSEC initialization"

    apk update
    apk upgrade
    apk add

    eend $?
}

stop() {
    ebegin ""
    eend $?
}
