# setup/packages/psec/psec-core/psec-core.post-install
#!/bin/sh

echo
echo --------------------------------------------
echo -------      PSEC Installation       -------
echo ---        script post-install.sh        ---
echo --------------------------------------------
echo

rc-service xenconsoled start
rc-service xenstored start
rc-service virtxend start

rc-update add start-domains
rc-service start-domains start
rc-update add psec-core-controller
rc-service psec-core-controller start

#/usr/lib/psec/bin/finish-core-init.sh

exit 0
