# setup/packages/psec/psec-core/psec-core.post-install
#!/bin/sh

echo
echo --------------------------------------------
echo -------      PSEC Installation       -------
echo ---        script post-install.sh        ---
echo --------------------------------------------
echo

rc-udpate add start-domains boot
rc-udpate add setup-pci boot
rc-update add attach-pci-devices boot
rc-update add xen-pci boot
rc-udpate add xenconsoled default
rc-update add xenstored default
rc-update add splash sysinit

/usr/lib/psec/bin/finish-core-init.sh
