# setup/packages/safecor/safecor-diag/diag.initd
#!/sbin/openrc-run

name="diag"
description="Diagnostic application for Safecor"
command="/usr/lib/safecor/diag/src/diag/start.sh"
pidfile="/var/run/$name.pid"
command_background=true

depend() {
    need xenstored udev-trigger debug
    after *
}
