# setup/packages/safecor/safecor-diag/safecor-diag.post-install
#!/bin/sh

#echo "Add debug service to default runlevel"
#rc-udpate add debug

#echo "Add diag app to default reunlevel"
#rc-update add diag

exit 0