# setup/packages/safecor/safecor-diag/safecor-diag.post-install
#!/bin/sh

rc-update add debug
rc-service debug start

exit 0